# root/etc/s6-overlay/s6-rc.d/init-nordlynx-iptables-rules/up
/etc/s6-overlay/s6-rc.d/init-nordlynx-iptables-rules/run
